# dpkg/prototype/DEBIAN/postinst
#!/bin/sh
update-rc.d ssh-tunnel defaults

