# dpkg/prototype/DEBIAN/prerm
#!/bin/sh
update-rc.d -f ssh-tunnel remove

